# systemd/paper-pulse-daily.timer
[Unit]
Description=Run Paper Pulse daily report generator every morning

[Timer]
OnCalendar=*-*-* 09:10:00
Persistent=true

[Install]
WantedBy=timers.target
