# systemd/paper-pulse-web.service
[Unit]
Description=Paper Pulse Web UI
After=network.target

[Service]
Type=simple
# TODO: change to your actual project path
WorkingDirectory=/root/paper-radar
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/root/paper-radar/.env
ExecStart=/root/paper-radar/.venv/bin/python -m uvicorn webapp:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
