# systemd/paper-pulse-daily.service
[Unit]
Description=Paper Pulse Daily Report
After=network.target

[Service]
Type=oneshot
# TODO: change to your actual project path
WorkingDirectory=/root/paper-radar
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/root/paper-radar/.env
ExecStart=/root/paper-radar/.venv/bin/python main.py
